# deploy/systemd/pelican-backup.timer
[Unit]
Description=Daily Pelican Raven recovery bundle backup
Documentation=file:/home/vinnieb58/projects/vulture/docs/current/PELiCAN_BACKUP.md

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=15m
Unit=pelican-backup.service

[Install]
WantedBy=timers.target
