# deploy/systemd/pelican-backup.service
[Unit]
Description=Pelican Raven recovery bundle backup (oneshot)
Documentation=file:/home/vinnieb58/projects/vulture/docs/current/PELiCAN_BACKUP.md

[Service]
Type=oneshot
User=vinnieb58
Group=vinnieb58
WorkingDirectory=/home/vinnieb58/projects/vulture
UMask=0077
ExecStart=/home/vinnieb58/projects/vulture/scripts/pelican_backup.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pelican-backup
NoNewPrivileges=true
TimeoutStartSec=3600
Restart=no
