# deploy/systemd/vulture-concert-watches.service
[Unit]
Description=Vulture concert watch cycle (scripts/run_concert_watches.py)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=vinnieb58
WorkingDirectory=/home/vinnieb58/projects/vulture
EnvironmentFile=/home/vinnieb58/projects/vulture/.env
ExecStart=/home/vinnieb58/projects/vulture/.venv/bin/python /home/vinnieb58/projects/vulture/scripts/run_concert_watches.py
