# deploy/systemd/vulture-scheduler.service
[Unit]
Description=Vulture hunt scheduler (single main.py cycle)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=vinnieb58
WorkingDirectory=/home/vinnieb58/projects/vulture
EnvironmentFile=/home/vinnieb58/projects/vulture/.env
ExecStart=/home/vinnieb58/projects/vulture/.venv/bin/python /home/vinnieb58/projects/vulture/main.py
